# config/vendor/etc/init.d/enforceperms
#!/sbin/openrc-run
# shellcheck shell=sh

description="Enforcing permissions on first boot"

start() {
    ebegin "Enforcing permissions on first boot"

    # Ensure /home/user is owned by `user`
    chown -R user:user /home/user
}

stop() {
    echo 'noop'
}
